# backend/supervisor_app.conf
[include]
files=/etc/supervisord.conf

[program:myapp]
directory=/root/cloud-disk/backend
command=gunicorn -w 4 -b 0.0.0.0:5000 app:app -t 3600 -c gunicorn_config.py

